# swift-app/com.mord58562.jury.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mord58562.jury.menubar</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/Applications/Jury.app/Contents/MacOS/jury-menubar</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/jury-menubar.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/jury-menubar.stderr.log</string>
</dict>
</plist>
